# blend-update.timer
[Unit]
Description=Scheduled BlendOS Update

[Timer]
OnCalendar=Thu *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target

