# blend-update.service
[Unit]
Description=Run BlendOS Akshara Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/blend-update
